# ops/systemd/ariadne-news-hourly.service
[Unit]
Description=Ariadne News Lane hourly ingest/compile
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/Users/damyandeshev/projects/ariadne
ExecStart=/usr/bin/env news-worker run-hourly
